# bridge/systemd/mt5-bridge.service
[Unit]
Description=MT5 AutoTrader bridge (Lovable dashboard)
After=network-online.target mt5linux.service
Wants=network-online.target
Requires=mt5linux.service

[Service]
Type=simple
User=__USER__
WorkingDirectory=__BRIDGE_DIR__
EnvironmentFile=__BRIDGE_DIR__/.env
ExecStart=__BRIDGE_DIR__/.venv/bin/python __BRIDGE_DIR__/main.py
Restart=always
RestartSec=5
StartLimitBurst=0

# Watchdog: if bridge stops writing HEALTH_FILE, systemd will restart it.
ExecStartPost=/bin/bash -c 'echo $MAINPID > /tmp/mt5-bridge.pid'

StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=multi-user.target